# NewsLetterService/crontab
0 0 * * * /usr/bin/python3 /app/collect.py >> /var/log/cron.log 2>&1
