# NewsLetterService/crontab
0 12 * * * python3 /app/_main.py >> /var/log/cron.log 
